# etc/mx-security/level1/files/mx-guard.service
[Unit]
Description=Moxa Security Utils Daemon Initial Setting

[Service]
Type=oneshot
ExecStart=/sbin/mx-guard start &
ExecStop=/sbin/mx-guard stop

